# scripts/systemd/hydra-autopilot.service
[Unit]
Description=Hydra Autopilot — long-running autonomous decision loop (overnight)
Documentation=https://github.com/gaberoo322/hydra/blob/master/docs/operator-playbooks/hydra-autopilot.md
After=network-online.target hydra-orchestrator.service
Wants=network-online.target

# Failure routing — reuse the existing notify-failure template.
OnFailure=hydra-notify-failure@%n.service

# Bounded auto-retry on failure. A transient Anthropic API 5xx killed the
# 2026-05-16 morning run after 1h 47min and the `Restart=no` policy left the
# system without an autopilot for ~10h until the next timer fire. One retry
# after 3 min absorbs transient errors; sustained failures still surface to
# the operator via OnFailure= once the burst cap is hit.
StartLimitBurst=2
StartLimitIntervalSec=3600

[Service]
Type=oneshot
WorkingDirectory=/home/gabe/hydra

# Soft cap above the autopilot's own 8h budget so Phase 7 (drain + digest) has
# headroom. Beyond this, systemd hard-kills.
RuntimeMaxSec=32400
TimeoutStopSec=600

# Subscription-billed: the claude CLI uses the operator's logged-in
# credentials in ~/.claude/. Do NOT set ANTHROPIC_API_KEY here.
Environment=HOME=/home/gabe
Environment=PATH=/home/gabe/.npm-global/bin:/home/gabe/.local/bin:/usr/local/bin:/usr/bin:/bin

# Optional per-run overrides (uncomment + edit if needed for a one-off):
# Environment=HYDRA_AUTOPILOT_TOKEN_BUDGET=2000000
# Environment=HYDRA_AUTOPILOT_MAX_SEC=28800
# Environment=HYDRA_AUTOPILOT_IDLE_TURNS=5

# The skill itself enforces all phase logic; the CLI just stays alive while the
# model produces tool calls and exits when the session terminates.
#
# --dangerously-skip-permissions is REQUIRED. Headless `claude -p` denies all
# confirmation-required tool calls by default; without skip-permissions, the
# subagents the autopilot dispatches return skill_denied and Phase 2 reaps
# them as immediate failures. Observed in the 2026-05-13 v1 smoke test.
ExecStart=/usr/bin/env claude --dangerously-skip-permissions -p "/hydra-autopilot"

StandardOutput=journal
StandardError=journal

# Auto-retry once on failure (bounded by StartLimitBurst above). The skill
# self-terminates cleanly on budget/clock/idle (exit 0 — no restart). Crashes
# (non-zero exit) get a single 3-min-delayed retry, which absorbs transient
# Anthropic API 5xx errors. Sustained failure hits the burst cap and surfaces
# to the operator via OnFailure=.
Restart=on-failure
RestartSec=180

[Install]
WantedBy=default.target
