# scripts/systemd/hydra-autopilot-morning.timer
[Unit]
Description=Hydra Autopilot — morning trigger (10:00 local time)
Documentation=https://github.com/gaberoo322/hydra/blob/master/docs/operator-playbooks/hydra-autopilot.md

[Timer]
# Local time. Pairs with hydra-autopilot.timer (22:00) so the autopilot fires
# twice per day. Both timers target the same hydra-autopilot.service oneshot;
# the service's 9h RuntimeMaxSec and the skill's internal 8h budget mean a
# morning run typically finishes by ~19:00, well before the evening fire.
OnCalendar=*-*-* 10:00:00

# If the machine was off at 10:00, run once at next boot. The skill's idle-drain
# logic ends the run quickly if there's nothing to do.
Persistent=true

# Tiny randomization so the autopilot doesn't fight other 10:00 timers for
# resources / API rate-limits.
RandomizedDelaySec=120

Unit=hydra-autopilot.service

[Install]
WantedBy=timers.target
