# scripts/systemd/hydra-dc-telemetry-snapshot.service
[Unit]
Description=Design-concept Phase B telemetry daily snapshot (issue #465)
Documentation=https://github.com/gaberoo322/hydra/issues/465

[Service]
Type=oneshot
Environment=HYDRA_ROOT=%h/hydra
WorkingDirectory=%h/hydra
ExecStart=/bin/bash -c '${HYDRA_ROOT}/scripts/autopilot/dc-telemetry-snapshot.sh'
# The script exits non-zero on hard failures (jq missing, telemetry
# endpoint unreachable, empty rollup). Let systemd log the failure;
# the next day's timer fires the snapshot fresh.
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
