# scripts/systemd/hydra-autopilot.service
[Unit]
Description=Hydra Autopilot — long-running autonomous decision loop
Documentation=https://github.com/gaberoo322/hydra/blob/master/docs/operator-playbooks/hydra-autopilot.md
After=network-online.target hydra-orchestrator.service
Wants=network-online.target

# Failure routing — reuse the existing notify-failure template.
OnFailure=hydra-notify-failure@%n.service

# Bounded auto-retry on failure. A transient Anthropic API 5xx killed the
# 2026-05-16 morning run after 1h 47min and the `Restart=no` policy left the
# system without an autopilot for ~10h until the next timer fire. One retry
# after 3 min absorbs transient errors; sustained failures still surface to
# the operator via OnFailure= once the burst cap is hit.
#
# Relaxed for issue #925: with StartLimitBurst=2 / IntervalSec=3600, two bad
# runs in an hour dead-zoned the unit for the rest of that hour — every
# subsequent Pace Gate launch (~every 15 min) was refused with "Start request
# repeated too quickly", starving the autopilot entirely. The dominant trigger
# was a `claude` self-exit code 143 (a child SIGTERM propagated as the parent's
# own exit STATUS) that SuccessExitStatus= did NOT cover (see below), so each
# clean self-exit armed Restart=on-failure and burned a start slot. With 143
# now declared a success, those exits no longer count against the burst; the
# wider burst (5) / shorter window (900s = one Pace Gate cycle) is belt-and-
# braces so a genuine transient burst cannot starve the Gate for more than a
# single cycle.
StartLimitBurst=5
StartLimitIntervalSec=900

[Service]
# Type=exec (NOT oneshot). This is a long-running foreground process — the
# claude CLI stays alive producing tool calls until the session ends — so exec
# is the correct type. It is ALSO load-bearing for RuntimeMaxSec below:
# systemd's RuntimeMaxSec has NO effect on Type=oneshot services (it only arms
# once a service enters the "running" substate, which oneshot never does), and
# oneshot's TimeoutStartSec defaults to infinity. Under the old `oneshot` a
# hung run (e.g. a subagent's `node --test` blocking on an open Redis handle)
# sat in `activating (start)` forever, never reaped, and blocked every
# subsequent timer fire. With Type=exec the 9h cap below actually fires.
Type=exec
WorkingDirectory=/home/gabe/hydra

# Hard cap above the autopilot's own 8h budget so Phase 7 (drain + digest) has
# headroom. Beyond this, systemd hard-kills the run (now effective — see the
# Type=exec note above; under oneshot this was silently ignored).
RuntimeMaxSec=32400
TimeoutStopSec=600

# A SIGTERM stop (`systemctl restart` / RuntimeMaxSec / shutdown) is a clean
# interrupt, not a failure (issue #898 / AC2). Without this, Type=exec reports
# the SIGTERM-killed unit as failed and the reap backstop sees
# EXIT_CODE=signal / EXIT_STATUS=TERM. The reap derivation maps that signal to
# `interrupted` regardless; declaring it a success here also keeps Restart=
# on-failure from spuriously retrying a deliberate restart.
#
# 143 / 130 added for issue #925. `SuccessExitStatus=SIGTERM` matches only a
# systemd-delivered *signal* (EXIT_CODE=signal / EXIT_STATUS=TERM); it does
# NOT match an exit *code* of 143. But the dominant termination here was the
# `claude` CLI self-exiting with code 143 (= 128+SIGTERM, the conventional code
# a Node/shell parent returns when a child it spawned dies on SIGTERM) — and
# 130 (= 128+SIGINT) for the Ctrl-C analogue. systemd sees those as
# EXIT_CODE=exited / EXIT_STATUS=143, marks the unit "Failed with result
# 'exit-code'", and arms Restart=on-failure → the StartLimit lockout above.
# Listing the numeric codes here declares those self-exits a success, so the
# unit no longer spuriously retries them and no longer burns start slots.
SuccessExitStatus=SIGTERM 143 130

# Subscription-billed: the claude CLI uses the operator's logged-in
# credentials in ~/.claude/. Do NOT set ANTHROPIC_API_KEY here.
Environment=HOME=/home/gabe
Environment=PATH=/home/gabe/.npm-global/bin:/home/gabe/.local/bin:/usr/local/bin:/usr/bin:/bin

# Optional per-run overrides (uncomment + edit if needed for a one-off):
# Environment=HYDRA_AUTOPILOT_TOKEN_BUDGET=2000000
# Environment=HYDRA_AUTOPILOT_MAX_SEC=28800
# Environment=HYDRA_AUTOPILOT_IDLE_TURNS=5

# The skill itself enforces all phase logic; the CLI just stays alive while the
# model produces tool calls and exits when the session terminates.
#
# ExecStart routes through the Pace Gate's exec mode (issue #1089 recurrence).
# Restart=on-failure relaunches this unit WITHOUT consulting the pace-gate
# timer's admission check, so a session-limit exit (code=1) used to storm:
# systemd relaunched into the still-exhausted quota every RestartSec=180s
# (journal 2026-06-10, restart counter 53→60) even though the #1112
# session-limit block was armed — that block only gated the timer path. The
# wrapper makes EVERY start (timer `systemctl start` AND systemd-internal
# Restart=) pass the same /api/usage/eligibility check: blocked/paused/
# fail-safe → clean exit 0 (Restart=on-failure disarms, no start slot
# burned); eligible → exec's the claude CLI below as the main process.
#
# The wrapped invocation is: claude --dangerously-skip-permissions -p
# "/hydra-autopilot". --dangerously-skip-permissions is REQUIRED. Headless
# `claude -p` denies all confirmation-required tool calls by default; without
# skip-permissions, the subagents the autopilot dispatches return
# skill_denied and Phase 2 reaps them as immediate failures. Observed in the
# 2026-05-13 v1 smoke test.
#
# %h/.local/bin/hydra-pace-gate.sh is the deploy-installed copy
# (scripts/deploy.sh keeps it AND this unit in sync with the repo on every
# merge to master).
ExecStart=/usr/bin/env bash %h/.local/bin/hydra-pace-gate.sh --exec-autopilot

# Reap-on-exit backstop (issue #898). Fires when the claude session exits for
# ANY reason that did not route through a term-check stop decision —
# print-mode session end, context exhaustion, SIGTERM from `systemctl restart`
# / RuntimeMaxSec, or a genuine crash. It POSTs run-end with a best-guess
# cause derived from systemd's $EXIT_CODE / $EXIT_STATUS (clean exit OR a
# SIGTERM/SIGINT signal-kill → `interrupted`; any other signal or non-zero
# exit → `crash`), so the run hash gets a terminal status
# exactly once instead of being left at `status: running` for the read-time
# sweeper to mislabel as `crash`. Idempotent: a no-op if term-check.py already
# POSTed run-end this turn. The `-` prefix makes a reap failure non-fatal to
# the unit stop.
ExecStopPost=-/usr/bin/env bash /home/gabe/hydra/scripts/autopilot/bootstrap.sh --reap

StandardOutput=journal
StandardError=journal

# Auto-retry once on failure (bounded by StartLimitBurst above). The skill
# self-terminates cleanly on budget/clock/idle (exit 0 — no restart). Crashes
# (non-zero exit) get a single 3-min-delayed retry, which absorbs transient
# Anthropic API 5xx errors. Sustained failure hits the burst cap and surfaces
# to the operator via OnFailure=.
#
# Issue #1089: the retry start runs the ExecStart wrapper above, so a relaunch
# into a blocked window (session-limit, operator pause, ahead-of-curve,
# eligibility unreachable) exits 0 cleanly instead of dying code=1 — the retry
# chain self-terminates after at most one cheap wrapper start.
Restart=on-failure
RestartSec=180

[Install]
WantedBy=default.target
